# test-gui-04/root/etc/skel/.config/autostart/setup_sforzato.desktop
[Desktop Entry]
Type=Application
Version=1.0
Name=setup_sforzato
Path=/usr/share/sforzato
Exec=setup_sforzato.sh
Icon=
Terminal=false
