# test-gui-03/root/etc/skel/.config/autostart/setup_mykiwi.desktop
[Desktop Entry]
Type=Application
Version=1.0
Name=setup_mykiwi
Path=/usr/share/mykiwi
Exec=setup_mykiwi.sh
Icon=
Terminal=false
